# prd-cron-backups/data/crontab
* * * * * root /data/tar_to_dropbox.sh
# An empty line is required at the end of this file for a valid cron file.
